# caddy-2/rootfs/etc/s6-overlay/s6-rc.d/caddy/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant App: Caddy 2
# Runs Caddy
# ==============================================================================

# Run App
exec /usr/bin/caddy.sh
