# caddy-2/rootfs/etc/services.d/caddy/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant  Add-on: Caddy 2
# Take down the S6 supervision tree when Caddy fails
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
  bashio::log.warning "Caddy 2 crashed, halting add-on"
  /run/s6/basedir/bin/halt
fi

bashio::log.info "Caddy 2 stoped, restarting..."